# backend/loop.service
[Unit]
Description=LOOP - Little Optical Output Pal
After=network.target
Wants=network.target

[Service]
Type=notify
User=pi
Group=pi
WorkingDirectory=/home/pi/loop/backend
Environment=PYTHONPATH=/home/pi/loop/backend
Environment=LOOP_ENV=production
ExecStart=/home/pi/loop/backend/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=3

# Permissions for GPIO and system control
SupplementaryGroups=gpio spi i2c dialout video audio plugdev netdev
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_NET_BIND_SERVICE

# Allow sudo for hotspot management (careful with this)
NoNewPrivileges=false

# Systemd integration
WatchdogSec=60
NotifyAccess=main

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loop

# Security (balanced with functionality)
PrivateTmp=true
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target 