# backend/loop.service
[Unit]
Description=LOOP - Little Optical Output Pal (Hardened)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=notify
User=${LOOP_USER}
Group=gpio
WorkingDirectory=${LOOP_PROJECT_DIR}/backend
Environment=PYTHONPATH=${LOOP_PROJECT_DIR}/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=${LOOP_PROJECT_DIR}/backend/venv/bin/python main.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=30
TimeoutStopSec=15
TimeoutSec=30

# Watchdog configuration - Pi-optimized
WatchdogSec=90
NotifyAccess=main

# Restart policy - aggressive but safe
Restart=always
RestartSec=10

# Process limits - Pi-optimized
LimitNOFILE=1024
LimitNPROC=64
MemoryMax=512M
CPUQuota=80%

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=${LOOP_PROJECT_DIR}/backend/media ${LOOP_PROJECT_DIR}/backend/config /home/${LOOP_USER}/.loop
DeviceAllow=/dev/spidev0.0 rw
DeviceAllow=/dev/gpiomem rw
DeviceAllow=char-gpio rw

# Allow binding to privileged ports (like port 80)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# GPIO access - minimal permissions
SupplementaryGroups=gpio spi

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loop

[Install]
WantedBy=multi-user.target 