# backend/loop.service
[Unit]
Description=LOOP - Little Optical Output Pal (Hardened)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=notify
User=pi
Group=gpio
WorkingDirectory=/home/pi/loop/backend
Environment=PYTHONPATH=/home/pi/loop/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/loop/backend/venv/bin/python main.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=30
TimeoutStopSec=15
TimeoutSec=30

# Watchdog configuration - Pi-optimized
WatchdogSec=90
NotifyAccess=main

# Restart policy - aggressive but safe
Restart=always
RestartSec=10

# Process limits - Pi-optimized
LimitNOFILE=1024
LimitNPROC=64
MemoryMax=512M
CPUQuota=80%

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/loop/backend/media /home/pi/loop/backend/config /home/pi/.loop
DeviceAllow=/dev/spidev0.0 rw
DeviceAllow=/dev/gpiomem rw
DeviceAllow=char-gpio rw

# GPIO access
SupplementaryGroups=gpio spi

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loop

[Install]
WantedBy=multi-user.target 