# backend/system-management.service
[Unit]
Description=Disable WiFi Power Management for LOOP
After=NetworkManager.service
Wants=NetworkManager.service
Before=loop.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=${LOOP_PROJECT_DIR}/backend/deployment/scripts/service-manager.sh setup-wifi

[Install]
WantedBy=multi-user.target 